# debian/cron.daily
#!/bin/sh
test -x /usr/sbin/automysqlbackup || exit 0

if [ ! -d /var/log/mysql ] ; then
    mkdir -p /var/log/mysql
    chown mysql:adm /var/log/mysql
    chmod 42750 /var/log/mysql
fi

/usr/sbin/automysqlbackup >> /var/log/mysql/automysqlbackup.log 2>&1
